# deploy/launchd/com.yumx.im-metaagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd 配置（macOS 进程守护）
  - RunAtLoad: 登录后自动启动
  - KeepAlive.Crashed: 崩溃（信号/非0退出）自动重启
  - KeepAlive.SuccessfulExit: 正常退出（exit 0）不重启（便于手动停机）
  - ThrottleInterval: 崩溃重启最小间隔，防止疯狂重启

  ⚠️ 路径占位符 __PROJECT_DIR__ 由 install.sh 自动替换为绝对路径。
     不要直接加载本文件，请运行 `bash deploy/install.sh`。
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.yumx.im-metaagent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__PROJECT_DIR__/deploy/start.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJECT_DIR__</string>

    <!-- 登录后自动启动 -->
    <key>RunAtLoad</key>
    <true/>

    <!-- 崩溃自动重启；正常退出(exit 0)不重启 -->
    <key>KeepAlive</key>
    <dict>
        <key>Crashed</key>
        <true/>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- 崩溃重启最小间隔（秒）-->
    <key>ThrottleInterval</key>
    <integer>15</integer>

    <!-- 日志输出 -->
    <key>StandardOutPath</key>
    <string>__PROJECT_DIR__/logs/im-metaagent.out.log</string>
    <key>StandardErrorPath</key>
    <string>__PROJECT_DIR__/logs/im-metaagent.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
